# docker/crontab
SHELL=/bin/bash
PATH=/app/venv/bin:/root/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
PYTHONPATH=/app/prism-insight

# daily maintenance
0 2 * * * cd /app/prism-insight && chmod +x tools/backup_configs.sh && ./tools/backup_configs.sh >> /app/prism-insight/logs/backup.log 2>&1
0 3 * * * cd /app/prism-insight && chmod +x tools/cleanup_logs.sh && ./tools/cleanup_logs.sh >> /app/prism-insight/logs/cleanup.log 2>&1
0 3 * * 0 cd /app/prism-insight && python3 compress_trading_memory.py >> /app/prism-insight/logs/compression.log 2>&1

# US pending order execution (Tue-Sat, KST)
5 10 * * 2-6 cd /app/prism-insight && python3 pending_order_batch.py >> /app/prism-insight/logs/pending_orders_$(date +\%Y\%m\%d).log 2>&1

# US market sessions (Tue-Sat, KST schedule)
15 0 * * 2-6 cd /app/prism-insight && python3 stock_analysis_orchestrator.py --mode morning >> /app/prism-insight/logs/morning_$(date +\%Y\%m\%d).log 2>&1
30 2 * * 2-6 cd /app/prism-insight && python3 stock_analysis_orchestrator.py --mode midday >> /app/prism-insight/logs/midday_$(date +\%Y\%m\%d).log 2>&1
30 6 * * 2-6 cd /app/prism-insight && python3 stock_analysis_orchestrator.py --mode afternoon >> /app/prism-insight/logs/afternoon_$(date +\%Y\%m\%d).log 2>&1
30 7 * * 2-6 cd /app/prism-insight && python3 performance_tracker_batch.py >> /app/prism-insight/logs/performance_$(date +\%Y\%m\%d).log 2>&1
0 8 * * 2-6 cd /app/prism-insight && python3 examples/generate_dashboard_json.py >> /app/prism-insight/logs/dashboard.log 2>&1
